# debian/prerm
#! /bin/sh

if [ -x /etc/init.d/sbs ]
then
    invoke-rc.d sbs stop || true
fi
if [ -x /sbin/systemctl ]
then
    systemctl stop sbs || true
fi
