# debian/prerm
#! /bin/sh

if [ -x /etc/init.d/sbs ]
then
	invoke-rc.d sbs stop || true
fi
